# resources/systemd/librephotos-frontend.service
[Unit]
Description=librephotos-frontend
[Service]
ExecStart=/usr/lib/librephotos/bin/librephotos-frontend
[Install]
WantedBy=multi-user.target
